# systemd/domum-media-backup.timer
[Unit]
Description=domum-media: schedule daily restic backup

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=15m
Persistent=true
Unit=domum-media-backup.service

[Install]
WantedBy=timers.target
