# systemd/domum-media-hot-prune.service
[Unit]
Description=domum-media: prune volatile hot storage
Documentation=file:///opt/domum-core-media/docs/SETUP-PLEX-BOOKS-HOT-STORAGE.md
ConditionPathExists=/usr/local/bin/domum-media

[Service]
Type=oneshot
ExecStart=/usr/local/bin/domum-media hot prune
TimeoutStartSec=1h
