# systemd/domum-media-check.service
[Unit]
Description=domum-media: restic repository integrity check
Documentation=file:///opt/domum-core-media/docs/SETUP-RESTIC.md
Wants=network-online.target
After=network-online.target
ConditionPathExists=/usr/local/bin/domum-media-backup

[Service]
Type=oneshot
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
# Weekly metadata check on all enabled targets, then run retention for any
# target configured with BACKUP_TARGET_*_RETENTION=weekly.
ExecStart=/usr/local/bin/domum-media-backup --check
ExecStart=/usr/local/bin/domum-media-backup --weekly-retention
TimeoutStartSec=4h
